# templates/default/profile.d/plenv.sh.erb
export PLENV_ROOT="/opt/plenv"
export PATH="$PLENV_ROOT/bin:$PATH"
eval "$(plenv init -)"
